# root/etc/services.d/bluetooth/run
#!/usr/bin/with-contenv bash

/usr/lib/bluetooth/bluetoothd -nd 
